# my-services/etc/systemd/user/mpd.socket.d/socket.conf
[Socket]
ListenStream=
ListenStream=%t/mpd.socket
ListenStream=6%U
