# my-services/lib/systemd/system/minissdpd.socket
[Unit]
Description=minissdpd socket

[Socket]
ListenStream=/var/run/minissdpd.sock

[Install]
WantedBy=sockets.target
