# my-services/lib/systemd/system/virtualbox.target
[Unit]
Description=Start virtualbox-services for users

[Install]
WantedBy=multi-user.target
