# my-services/usr/lib/systemd/user/hplip-systray.service
[Unit]
Description=HPLIP systemtray	

[Service]
ExecStart=/usr/bin/hp-systray -x
Restart=on-failure

[Install]
WantedBy=default.target
