# my-services/etc/systemd/system/user@1002.service.d/dir.conf
[Unit]
Requires=home-gast.service
After=home-gast.service

[Service]
ExecStopPost=/bin/rm -rf /home/gast
ExecStartPre=/bin/echo BLAAA %u %U
