# my-services/etc/systemd/user/mpd.service.d/service.conf
[Unit]
ConditionPathExists=%h/.config/mpd/mpd.conf
